# debian/DEBIAN/postinst
#!/bin/bash

systemctl --no-reload enable auto_update_miniboot.service
exit 0

